# config/supervisord_web.conf
[supervisord]

[supervisorctl]

[program:web-app]
command = /usr/local/bin/uwsgi --ini /app/config/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx